# System_services/gui.service
[Unit]
Description=Boat GUI Display
After=canup.service

[Service]
user=current
WorkingDirectory=/home/current/github/current-app
Environment=Display=:0
Environment=XAUTHORITY=/home/current/.Xauthority
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/local/bin/pipenv run python app/app.py
Restart=always
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target
